# dashboard/launchd/com.haoyining.dashboard-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.haoyining.dashboard-sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/haoyining/.openclaw/workspace/dashboard/sync_to_repo.sh</string>
    <string>/Users/haoyining/.openclaw/workspace</string>
    <string>new-看板</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>10</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/haoyining/.openclaw/workspace/dashboard/launchd/dashboard-sync.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/haoyining/.openclaw/workspace/dashboard/launchd/dashboard-sync.err.log</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
